# Linux/Dynamic wallpaper/dynamic_wallpaper.service
[Unit]
Description=Set dynamic wallpaper

[Service]
Type=oneshot
ExecStart=/home/v/bin/dynamic_wallpaper.py
