# Linux/Config files/.config/systemd/user/notify_9_30.service
[Unit]
Description=Desktop notification at 9:30 PM

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "9:30" "Last graceful call"
