# Linux/Config files/.config/systemd/user/notify_9.service
[Unit]
Description=Desktop notification at 9:00 PM

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "9:00" "Tie loose ends\nD B+1"
