# Linux/Config files/etc/systemd/system/nightly-shutdown.service
[Unit]
Description=Nightly Shutdown

[Service]
Type=oneshot
ExecStart=/usr/bin/shutdown now

