# Linux/Config files/etc/systemd/system/primary-monitor.service
[Unit]
Description= Make DP-1 primary

[Service]
Type=simple
ExecStart= /usr/bin/monitor-primary.sh

[Install]
WantedBy=multi-user.target
